# docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:php]
command=php-fpm
autostart=true
autorestart=true
